# pi_temp_reader_mqtt.service
[Unit]
Description=Python temp reader mqtt
After=multi-user.target

[Service]
Type = idle
ExecStart=/usr/bin/python3 -u /home/pi/git/dev/rasp_pi_heater_read/temp_reader_mqtt.py
WorkingDirectory=/home/pi/git/dev/rasp_pi_heater_read
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi

[Install]
WantedBy=multi-user.target
